# roles/gnome/files/autostart/1password.desktop
[Desktop Entry]
Name=1Password
Exec=/home/linuxbrew/.linuxbrew/bin/1password --silent
Terminal=false
Type=Application
Icon=1password
StartupWMClass=1Password
Comment=Password manager and secure wallet
MimeType=x-scheme-handler/onepassword;x-scheme-handler/onepassword8;
Categories=Office;
X-GNOME-Autostart-enabled=true
